# docker/s6-overlay/s6-overlay/s6-rc.d/victoriametrics/run
#!/command/with-contenv sh
exec /usr/local/bin/victoria-metrics \
    --storageDataPath=/data/vm \
    --retentionPeriod=90d \
    --httpListenAddr=:8428 \
    --influxSkipSingleField
